# infra/deploy-server/deploy-server.service
[Unit]
Description=LifeSync360 Deploy Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/deploy-server
ExecStart=/opt/deploy-server/venv/bin/python deploy_server.py
Restart=always
RestartSec=5
EnvironmentFile=/etc/deploy-server/env
Environment=ANSIBLE_DIR=/opt/ansible/onprem-prod-repo/ansible
Environment=VAULT_PASS_FILE=/home/ubuntu/.vault_pass
Environment=DEPLOY_LOG=/var/log/ansible-deploy.log

[Install]
WantedBy=multi-user.target
