# onprem-prod-repo/ansible/roles/private_api/templates/private-api.service.j2
[Unit]
Description=LifeSync360 Private API
After=network.target

[Service]
Type=simple
User=ansible
WorkingDirectory=/opt/private-api
ExecStart=/opt/private-api/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5
Environment=ENV=production
Environment=ANSIBLE_ENV=production
Environment=AWS_DEFAULT_REGION=ap-northeast-2
Environment=DEPLOY_TOKEN={{ deploy_token }}
Environment=CONTROL_NODE_URL={{ control_node_url }}
Environment=TOKEN_AES_KEY_B64={{ vault_token_aes_key_b64 }}
Environment=DB_HOST={{ mysql_vm_ip }}
Environment=DB_USER={{ mysql_app_user }}
Environment=DB_PASS={{ mysql_app_password }}

[Install]
WantedBy=multi-user.target
